# crontab
SHELL=/bin/bash
0,15,30,45 * * * * root /reddit-download.sh > /var/log/cron.log 2>&1
# An empty line is required at the end of this file for ka valid cron file.
