# crontab
SHELL=/bin/bash
0/5 * * * * root /reddit-download.sh
# An empty line is required at the end of this file for ka valid cron file.